# supervisor.cnpmjs.conf
[program:cnpmjs]
command=/usr/bin/node --harmony /usr/local/share/cnpmjs/dispatch.js 
